# systemd/tg-digest-1100.timer
[Unit]
Description=Schedule Web3 Digest at 11:00

[Timer]
OnCalendar=*-*-* 11:00:00 Asia/Jerusalem
Persistent=true

[Install]
WantedBy=timers.target
