# systemd/tg-digest-1100.service
[Unit]
Description=Send Web3 Digest at 11:00 (Asia/Jerusalem)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/tg-web3-digest
EnvironmentFile=/opt/tg-web3-digest/.env
ExecStart=/opt/tg-web3-digest/.venv/bin/python /opt/tg-web3-digest/tg_digest_sender.py
